# scripts/cron/heartbeat-betterstack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ecypro.heartbeat.betterstack</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>scripts/heartbeat.mjs</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/tmp/ecypro-heartbeat.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ecypro-heartbeat.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
